# src/svxlink/systemd/svxlink.service.in

[Unit]
Description=SvxLink repeater control software
Documentation=man:svxlink(1) man:svxlink.conf(5)

[Service]
Type=simple
Environment=RUNASUSER=@SVXLINK_USER@
Environment=STATEDIR=@SVX_LOCAL_STATE_DIR@
Environment=CFGFILE=@SVX_SYSCONF_INSTALL_DIR@/svxlink.conf
Environment=LOGFILE=@LOCAL_STATE_DIR@/log/svxlink.log
EnvironmentFile=@SYSTEMD_DEFAULTS_FILES_DIR@/svxlink
UMask=0002
ExecStartPre=/bin/sh -c ' \
	if echo -n "${LOGFILE}" | grep -q '^/'; then \
	  touch "${LOGFILE}"; \
	  chmod u+w "${LOGFILE}"; \
	  chown ${RUNASUSER}:$(id -gn ${RUNASUSER}) "${LOGFILE}"; \
	fi; \
	'
ExecStartPre=/bin/sh -c ' \
	mkdir -m0775 -p "${STATEDIR}"; \
	chmod -R u+rwX "${STATEDIR}"; \
	chown -R ${RUNASUSER}:$(id -gn ${RUNASUSER}) "${STATEDIR}"; \
	'
ExecStart=@BIN_INSTALL_DIR@/svxlink --logfile=${LOGFILE} --config=${CFGFILE} --runasuser=${RUNASUSER}
ExecStopPost=@BIN_INSTALL_DIR@/svxlink --config=${CFGFILE} --runasuser=${RUNASUSER} --reset --quiet
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
Nice=-10
TimeoutStartSec=60
TimeoutStopSec=10
LimitCORE=infinity
WorkingDirectory=@SVX_SYSCONF_INSTALL_DIR@

[Install]
WantedBy=multi-user.target
